# auto_update.service
[Unit]
Description=Auto update biotron and jig firmware

[Service]
ExecStart=/bin/bash ~/initial_script.sh

[Install]
WantedBy=multi-user.target